# pioneer_ws/scripts/ros2_start.service
[Unit]
Description=ROS 2 Enhanced Auto Start for %i
After=network-online.target
Wants=network-online.target
# Wait for system to be fully ready
After=graphical-session.target

[Service]
Type=simple
# Environment variables
Environment="ROBOT_ID=p4"
Environment="ROS_DOMAIN_ID=0"
Environment="ROS_LOCALHOST_ONLY=0"
Environment="DISPLAY=:0"
Environment="HOME=/home/pioneer-iv"

# User and working directory
User=pioneer-iv
Group=pioneer-iv
WorkingDirectory=/home/pioneer-iv/ros2_ws

# Enhanced startup script
ExecStart=/bin/bash /home/pioneer-iv/ros2_ws/pioneer_ws/scripts/start_ros2.sh

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security settings
NoNewPrivileges=true
# Allow access to devices
SupplementaryGroups=dialout i2c gpio

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ros2-startup

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStartSec=120
TimeoutStopSec=30

# Health check (optional)
ExecStartPost=/bin/bash -c 'sleep 30 && /home/pioneer-iv/ros2_ws/pioneer_ws/scripts/health_check.sh'

[Install]
WantedBy=multi-user.target